# nodechecker/nodechecker/etc/profile.d/oi.sh.x86_64
export OI_HOST_GROUP=oi-group
export OI_ROOT=/opt/openinfinity/3.1.0
export OI_HEALTH_MONITORING_ROOT=$OI_ROOT/healthmonitoring/nodechecker
export OI_RRD_HTTP_SERVER_ROOT=$OI_ROOT/healthmonitoring/rrd-http-server
export OI_COLLECTD_ROOT=$OI_ROOT/healthmonitoring/collectd
export OI_CLOUD_ZONE=oi-zone
export OI_INSTANCE_ID=1
export OI_CLUSTER_ID=1
export OI_MACHINE_ID=1
export OI_PUBLIC_IP=1.2.3.4
export JAVA_HOME=/etc/alternatives/jre_openjdk/

